# ops/crontab
0 0 * * * /srv/vzdv/vzdv-import
0 * * * * /srv/vzdv/backup.sh
